# roles/sqoop/templates/etc/profile.d/sqoop.sh.j2
# {{ ansible_managed }}
# Apache Sqoop

export SQOOP_HOME={{ sqoop_source_link }}
export PATH=$PATH:$SQOOP_HOME/bin
